# board/bunzo/common/rootfs-overlay/etc/systemd/system/bunzo-shell.service
[Unit]
Description=bunzo chat shell (TUI stub)
After=systemd-user-sessions.service
ConditionKernelCommandLine=!bunzo.recovery

[Service]
Type=simple
ExecStartPre=/bin/sleep 3
ExecStart=/usr/bin/bunzo-shell
Restart=always
RestartSec=1
StandardInput=tty
StandardOutput=tty
StandardError=tty
TTYPath=/dev/ttyAMA0
TTYReset=yes
TTYVHangup=yes
UtmpIdentifier=bunzo
UtmpMode=user
Environment=TERM=linux
Environment=COLUMNS=80
Environment=LINES=24
Environment=BUNZO_SHELL_MODE=serial

[Install]
WantedBy=multi-user.target
